# gluster/services.d/gluster-exporter/run
#!/usr/bin/execlineb -P
/sbin/gluster-exporter
