# gluster/services.d/udevd/run
#!/usr/bin/execlineb -P
/lib/systemd/systemd-udevd
