# gluster/services.d/glustereventsd/run
#!/usr/bin/execlineb -P
/sbin/glustereventsd
